# deploy/go_zoo.service
[Unit]
Description=go_zoo

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=/root/go_zoo/zoo

[Install]
WantedBy=multi-user.target
